# deploy/roles/nginx/files/nginxservice.service
[Unit]
Description=Systemd service for nginx.

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/nginx.sh

[Install]
WantedBy=multi-user.target